# test/copy.conf
[program:copy-extension]
command=/opt/copy-extension.sh
startsecs=0
autorestart=false
priority=100
stdout_logfile=/var/log/neko/copy.log
redirect_stderr=true
